# infrastructure/templates/salem-engine.service.j2
[Unit]
Description=Salem 1692 Engine
After=network.target postgresql.service

[Service]
Type=simple
User={{ app_user }}
WorkingDirectory={{ project_root }}
EnvironmentFile=/etc/salem-1692/env
ExecStart={{ project_root }}/engine/salem
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
